# deployment/cridergpt-engine-update.timer
[Unit]
Description=Check GitHub for CriderGPT Engine updates every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=cridergpt-engine-update.service

[Install]
WantedBy=timers.target
